# weetee.service
[Unit]
Description=Phone home if IP has changed

[Service]
User=root
ExecStart=/usr/local/bin/weetee

[Install]
WantedBy=basic.target

